# docker/crontab
0 0 * * * root /usr/local/bin/node /app/run.js >> /var/log/cron.log 2>&1
