# docker/crontab
# Cron job that sends emails every day
0 0 * * * root /bin/bash /app/docker/runner.bash >> /var/log/cron.log 2>&1
